# deploy/systemd/kiana-verify-chain.timer
[Unit]
Description=Run KI_ana verify_chain periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

