# deploy/systemd/kiana-knowledge-prune.timer
[Unit]
Description=Run Knowledge Prune daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=kiana-knowledge-prune.service

[Install]
WantedBy=timers.target
