# deploy/systemd/kiana-device-events-stats.service
[Unit]
Description=KI_ana Device Events Stats Rollup (hourly)
After=network.target

[Service]
User=kiana
WorkingDirectory=/home/kiana/ki_ana
EnvironmentFile=/home/kiana/ki_ana/.env
ExecStart=/bin/bash -lc 'source /home/kiana/ki_ana/.venv/bin/activate 2>/dev/null || true; python3 system/device_events_stats_rollup.py'
Restart=no

[Install]
WantedBy=multi-user.target
